# startup.service
[Unit]
Description=Digital Bulletin Board Display Service
Documentation=https://github.com/jklul47/digi_bulletin
After=graphical-session.target multi-user.target network-online.target
Wants=graphical-session.target
Requires=network-online.target

[Service]
Type=simple
User=nljk
Group=nljk
WorkingDirectory=/home/nljk/dev/digi_bulletin

# Environment variables for GUI applications
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/nljk/.Xauthority
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=HOME=/home/nljk
Environment=USER=nljk
Environment=LOGNAME=nljk
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Execute the startup script
ExecStart=/bin/bash /home/nljk/dev/digi_bulletin/startup.sh

# Restart policy - be more conservative
Restart=on-failure
RestartSec=30
StartLimitBurst=2
StartLimitIntervalSec=600

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bulletin-board

# Security settings (relaxed for GUI apps)
NoNewPrivileges=false
PrivateTmp=false

# Timeout settings - longer for GUI startup
TimeoutStartSec=120
TimeoutStopSec=30

[Install]
WantedBy=graphical.target